# docker/php_cli/crontab
*/5 * * * * cd /var/www/html/nt && php aci/index.php cli/Tracking in_resque
